# l8.service
[Unit]
Description=L8 proxy delay server

[Service]
Type=simple
ExecStart=/usr/local/bin/l8

[Install]
WantedBy=multi-user.target
